# dist/conf/system/gofsync.service
[Unit]
Description=golang foreman synchronization web api

[Service]
PIDFile=/tmp/gofsync.pid-8086
User=root
Group=root
WorkingDirectory=/opt/goFsync
ExecStart=/bin/bash -c '/opt/goFsync/fGetter -server -file=/opt/goFsync/webHosts'
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
#ExecReload="/usr/bin/pkill fGetter;/bin/bash -c '/opt/goFsync/fGetter'"
#ExecStop=/usr/bin/pkill fGetter

[Install]
WantedBy=multi-user.target
